# files/system/etc/xdg/autostart/steam.desktop
[Desktop Entry]
Type=Application
Name=Steam (Disabled)
Exec=/usr/bin/steam
Hidden=true
NoDisplay=true
X-GNOME-Autostart-enabled=false
